# simulation_ws/devel/etc/catkin/profile.d/25.turtlebot-stdr.sh
/home/comp2/simulation_ws/devel/.private/turtlebot_stdr/etc/catkin/profile.d/25.turtlebot-stdr.sh