# simulation_ws/devel/etc/catkin/profile.d/25.turtlebot-navigation.sh
/home/comp2/simulation_ws/devel/.private/turtlebot_navigation/etc/catkin/profile.d/25.turtlebot-navigation.sh